# packer/first-run/sensor-provision.service
[Unit]
Description=Honeypot Sensor First-Boot Provisioning
Documentation=https://github.com/your-org/honeypot-pr
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=oneshot
ExecStart=/opt/sensor/sensor-provision.sh
RemainAfterExit=no
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
